# mtp.mount
[Mount]
Type=mtp
Exec=/usr/lib/gvfsd-mtp
AutoMount=false
